# glance-api.service
[Unit]
Description=glance API server

[Service]
ExecStart=/usr/bin/glance-api
User=glance

[Install]
WantedBy=multi-user.target
